# seed/data/etc/init/git-srv.conf
description "Start/Stop Git fetch server"
start on startup
stop on shutdown
exec /bin/bash -c "/usr/bin/git daemon --user=git-srv --group=git-srv --verbose --export-all --reuseaddr --base-path=/srv/git/public /srv/git/public &> /var/log/git-srv.log"
respawn
