# seed/data/etc/cron.d/upgrade
0 5	* * *	root	DEBIAN_FRONTEND=noninteractive apt-get update && apt-get -q -y dist-upgrade && apt-get clean && apt-get autoremove
