# seed/data/etc/cron.d/upgrade
PATH=/usr/bin:/bin:/usr/sbin:/sbin
DEBIAN_FRONTEND=noninteractive
0 5	* * *	root	apt-get update && apt-get -q -y dist-upgrade && apt-get clean && apt-get autoremove
