# src/certbot/~data/etc/cron.daily/certbot-renew
#!/bin/sh
PATH=/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin
certbot renew -n %POST_HOOK%>> /var/log/certbot-renew.log 2>&1
