# src/data/etc/init/hosts-fix.conf
# hosts - fix system hosts
#
# This task is run on startup to fix the system hosts in /etc/hosts,
# fixing localhost.localdomain because it is failing when multilpe
# aliasses are registered for one IP per line.

description "Fix system hosts"

start on startup

task
exec cat /etc/hosts | sed 's/localhost.localdomain localhost/localhost localhost.localdomain/g' | tee /etc/hosts > /dev/null
