# src/data/etc/cron.d/upgrade
PATH=/usr/bin:/bin:/usr/sbin:/sbin
DEBIAN_FRONTEND=noninteractive
0 5 * * * root ( apt-get -q -y dist-upgrade && apt-get -y autoremove && apt-get clean ) >> /var/log/upgrade.log 2>&1
